# ros.service
[Unit]
Description=ros platform service
After=network.target
Wants=network.target

[Service]
User=szymonri
ExecStart=SCRIPT_FULL_PATH/build_and_run.sh

[Install]
WantedBy=multi-user.target


